# systemd/wg-admin.service
[Unit]
Description=WireGuard Admin Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/wg-manager
Environment="PYTHONPATH=/opt/wg-manager"
EnvironmentFile=/opt/wg-manager/.env
ExecStart=/opt/wg-manager/venv/bin/uvicorn backend.admin.main:app --host ${ADMIN_HOST:-127.0.0.1} --port ${ADMIN_PORT:-8081}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
